# contrib/init/bitcoind.conf
description "Lyokocoin ABC Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env LYOKOCOIND_BIN="/usr/bin/lyokocoind"
env LYOKOCOIND_USER="lyokocoin"
env LYOKOCOIND_GROUP="lyokocoin"
env LYOKOCOIND_PIDDIR="/var/run/lyokocoind"
# upstart can't handle variables constructed with other variables
env LYOKOCOIND_PIDFILE="/var/run/lyokocoind/lyokocoind.pid"
env LYOKOCOIND_CONFIGFILE="/etc/lyokocoin/lyokocoin.conf"
env LYOKOCOIND_DATADIR="/var/lib/lyokocoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # lyokocoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$LYOKOCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run lyokocoind."
        echo "The setting must appear in $LYOKOCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Lyokocoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$LYOKOCOIND_PIDDIR"
    chmod 0755 "$LYOKOCOIND_PIDDIR"
    chown $LYOKOCOIND_USER:$LYOKOCOIND_GROUP "$LYOKOCOIND_PIDDIR"
    chown $LYOKOCOIND_USER:$LYOKOCOIND_GROUP "$LYOKOCOIND_CONFIGFILE"
    chmod 0660 "$LYOKOCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$LYOKOCOIND_PIDFILE" \
    --chuid $LYOKOCOIND_USER:$LYOKOCOIND_GROUP \
    --exec "$LYOKOCOIND_BIN" \
    -- \
    -pid="$LYOKOCOIND_PIDFILE" \
    -conf="$LYOKOCOIND_CONFIGFILE" \
    -datadir="$LYOKOCOIND_DATADIR" \
    -disablewallet \
    -daemon

